# ax3600/files/etc/uci-defaults/99-custom
uci -q batch << EOF
set wireless.radio0.disabled='1'
set wireless.radio0.country='US'
set wireless.radio1.mu_beamformer='1'
set wireless.radio1.country='US'
set wireless.radio2.channel='auto'
set wireless.radio2.mu_beamformer='1'
set wireless.radio2.country='US'
commit wireless
EOF

uci -q batch << EOF
set network.lan.ipaddr='192.168.31.1'
commit network
set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-iface[0].ssid='ImmortalWrt'
commit wireless
set dhcp.@dnsmasq[0].rebind_protection="0"
commit dhcp
EOF
