# x86/files/etc/uci-defaults/99-custom
uci -q batch << EOF
set dhcp.@dnsmasq[0].rebind_protection="0"
commit dhcp
EOF
