# rockchip/files/etc/uci-defaults/99-custom
uci -q batch << EOF
set network.lan.ipaddr='192.168.31.1'
commit network
set dhcp.@dnsmasq[0].rebind_protection="0"
commit dhcp
EOF
